# deploy/systemd/hermes-ui.service.example
[Unit]
Description=NewsCraft AI SvelteKit UI
After=network-online.target newsroom-harness.service
Wants=network-online.target newsroom-harness.service

[Service]
Type=simple
WorkingDirectory=/home/jigar/hermes-ui
EnvironmentFile=/home/jigar/hermes-ui/.env
Environment=NODE_ENV=production
Environment=HOST=127.0.0.1
Environment=PORT=3001
ExecStart=/usr/bin/env node build/index.js
Restart=on-failure
RestartSec=5
KillSignal=SIGTERM
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
